# volume/etc/systemd/system/mongod.service
[Unit]
Description=MongoDB server
After=network.target

[Service]
ExecStart=mongod
Restart=always

[Install]
WantedBy=default.target
