# volume/etc/supervisor/conf.d/mongod.conf
[program:mongod]
command=mongod
autostart=true
autorestart=true
stderr_logfile=/var/log/sup-mongod.err.log
stdout_logfile=/var/log/sup-mongod.out.log
