# volume/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/sup-nginx.err.log
stdout_logfile=/var/log/sup-nginx.out.log
